# services/backend
#!/bin/sh

# PROVIDE: backend
# REQUIRE: postgresql ldap
# BEFORE: nginx

. /etc/rc.subr

name=backend
desc="backend jail"
