# services/ldap
#!/bin/sh

# PROVIDE: ldap
# REQUIRE: letsencrypt

. /etc/rc.subr

name=ldap
desc="ldap jail"
