# services/letsencrypt
#!/bin/sh

# PROVIDE: letsencrypt

. /etc/rc.subr

name=letsencrypt
desc="letsencrypt jail"
