# services/nginx
#!/bin/sh

# PROVIDE: nginx
# REQUIRE: letsencrypt

. /etc/rc.subr

name=nginx
desc="nginx jail"
