# services/postgresql
#!/bin/sh

# PROVIDE: postgresql
# BEFORE: nginx

. /etc/rc.subr

name=postgresql
desc="postgresql jail"
